# systemd/altlinker.service.d/override.conf
[Service]
Environment="PYTHONPATH=/home/folan/altlinker/.venv/bin/python3.10"
Environment="TELEGRAM_BOT_TOKEN=token"
